# docker/i2b2-pgsql/supervisord.conf
[supervisord]
nodaemon=true

[program:postgresql]
command=/etc/init.d/postgresql start

autostart=true
autorestart=true